# vg-vnc/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24
autostart=true
autorestart=true
priority=1

[program:fluxbox]
command=/usr/bin/fluxbox
autostart=true
autorestart=true
priority=2

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -forever -shared -nopw
autostart=true
autorestart=true
priority=3

[program:novnc]
command=/opt/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080
autostart=true
autorestart=true
priority=4

[program:chromium]
command=/usr/bin/chromium-browser --no-sandbox --disable-dev-shm-usage --remote-debugging-port=9222 --disable-gpu --display=:1
autostart=true
autorestart=true
priority=5

[program:nodejs]
command=npm run start:prod
directory=/app
autostart=true
autorestart=true
priority=6
stdout_logfile=/var/log/supervisor/nodejs.log
stderr_logfile=/var/log/supervisor/nodejs.err

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=7
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
